# setup/etc/rc.d/rc.local
/usr/sbin/inetd &
# start qconn for remote connections with Momentics
/usr/sbin/qconn &
# start SSH daemon
/usr/sbin/sshd &